# systemd/constitutional-ai-backend.service
[Unit]
Description=Constitutional AI Backend
After=network.target constitutional-ai-llm.service
Wants=constitutional-ai-llm.service

[Service]
Type=simple
WorkingDirectory=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/backend
Environment="PATH=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="RAG_SIMILARITY_THRESHOLD=0.01"
ExecStart=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/backend/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8900
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=default.target
