# apps/constitutional-dashboard/constitutional-cast.service
[Unit]
Description=Constitutional Dashboard Auto-Cast to Sovis
After=network.target constitutional-dashboard.service
Requires=constitutional-dashboard.service

[Service]
Type=simple
User=ai-server
WorkingDirectory=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/apps/constitutional-dashboard
Environment="PATH=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/venv/bin:/home/ai-server/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CAST_DISPLAY_NAME=Sovis"
Environment="DASHBOARD_PUBLIC_URL=http://192.168.86.32:5000"
Environment="CAST_KEEPALIVE_SECONDS=30.0"
Environment="CAST_RECAST_BACKOFF_SECONDS=5.0"
Environment="DASHBOARD_STARTUP_WAIT_SECONDS=20.0"
ExecStart=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/venv/bin/python3 /home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/apps/constitutional-dashboard/cast_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
