# apps/constitutional-dashboard/constitutional-dashboard.service
[Unit]
Description=Constitutional Dashboard (Sovis)
After=network.target

[Service]
Type=simple
User=ai-server
WorkingDirectory=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/apps/constitutional-dashboard
Environment="PATH=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DASHBOARD_HOST=0.0.0.0"
Environment="DASHBOARD_PORT=5000"
Environment="DASHBOARD_PUBLIC_URL=http://192.168.86.32:5000"
Environment="CAST_DISPLAY_NAME=Sovis"
ExecStart=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/venv/bin/python3 /home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/apps/constitutional-dashboard/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
